# files/image_config/rps/configure-rps@.service
[Unit]
Description=Configure RPS for %i
After=warmboot-finalizer.service

[Service]
Type=oneshot
ExecStart=/usr/bin/rps.py %i

